# ansible/roles/backup/templates/kopia-auto-backup.timer.j2
[Unit]
Description=Automatic kopia backup

[Timer]
OnCalendar={{ kopia_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
